# docs/snippets/ssh/conf/mysqld.service.d/custom.conf
[Service]
LimitNOFILE = infinity
LimitMEMLOCK = infinity
LimitSTACK = infinity